# mcs-matrix-tx.service
[Unit]
Description=MCS Matrix TX – mesh network + throughput sweep
After=network-online.target multi-user.target systemd-modules-load.service
Wants=network-online.target
After=systemd-networkd-wait-online.service NetworkManager-wait-online.service

[Service]
Type=oneshot

# --- Configuration (edit these) ---
Environment=TX_BAT_IP=192.168.50.1/24
Environment=RX_IP=192.168.50.2
Environment=CTRL_PORT=9999
Environment=IPERF_PORT=5201
Environment=IPERF_DURATION=10
Environment=INITIAL_WAIT=0
Environment=GUARD=0
Environment=MESH_SCRIPT=/opt/mcs-test/sdmah-mesh.sh
Environment=TX_SCRIPT=/opt/mcs-test/tx_matrix.py

# Step 0: Stop lighttpd and NetworkManager
ExecStartPre=-/bin/systemctl stop lighttpd
ExecStartPre=-/bin/systemctl stop NetworkManager

# Step 0b: Bring up wlan1 AP interface and restart hostapd
ExecStartPre=-/sbin/ip link set wlan1 up
ExecStartPre=-/bin/systemctl restart hostapd

# Step 1: Load morse driver modules
ExecStartPre=/sbin/modprobe cfg80211
ExecStartPre=/sbin/modprobe mac80211
ExecStartPre=/sbin/modprobe crc7
ExecStartPre=-/sbin/rmmod morse dot11ah
ExecStartPre=/sbin/insmod /home/pi/morse_driver/dot11ah/dot11ah.ko
ExecStartPre=/sbin/insmod /home/pi/morse_driver/morse.ko

# Step 1: Make sure WiFi radio is unblocked
ExecStartPre=/usr/sbin/rfkill unblock wifi

# Step 2: Patch the BAT_IP in sdmah-mesh.sh and run it
ExecStartPre=/bin/bash -c 'sed -i "s|^BAT_IP=.*|BAT_IP=\"${TX_BAT_IP}\"|" ${MESH_SCRIPT}'
ExecStartPre=/bin/bash ${MESH_SCRIPT}

# Step 4: Small settle delay for the mesh to form
ExecStartPre=/bin/sleep 5

# Step 6: Run the MCS sweep
ExecStart=/bin/bash -c 'exec env INITIAL_WAIT=${INITIAL_WAIT} GUARD=${GUARD} python3 ${TX_SCRIPT} ${RX_IP} ${CTRL_PORT} ${IPERF_PORT} ${IPERF_DURATION}'

# Step 7: Start drone server after sweep completes
ExecStartPost=/bin/systemctl start drone-server

WorkingDirectory=/opt/mcs-test
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
